# visionect_server_aio/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash
echo "Starting PostgreSQL server..."
exec su - postgres -c "postgres -D /data/postgres"
