# visionect_server/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash
su - postgres -c "/usr/bin/postgres -D /data/postgres"
