# visionect_server/rootfs/etc/services.d/redis/run
#!/usr/bin/env bash
set -e

# Function to log messages
log() {
    echo "[redis] $1"
}

log "Starting Redis server..."

# Get Redis port from environment
REDIS_PORT="${REDIS_PORT:-6379}"

log "Redis will listen on port: $REDIS_PORT"

# Start Redis server
exec redis-server \
    --port "$REDIS_PORT" \
    --bind 127.0.0.1 \
    --save "" \
    --appendonly no \
    --protected-mode no \
    --loglevel warning
