# visionect_server/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bash
exec redis-server
